# backup.timer
[Unit]
Description=Backup Timer

[Timer]
OnCalendar=*-*-* 00:15:30

[Intall]
WantedBy=timers.target
